# files/system/sneexy/wuzetka/etc/skel/.config/autostart/unfuck_music.desktop
[Desktop Entry]
Exec=/usr/bin/unfuck_music
Icon=application-x-shellscript
Name=unfuck_music
Type=Application
X-KDE-AutostartScript=true
