# files/system/sneexy/wuzetka/etc/profile.d/wuzetka-settings-xdg.sh
export XDG_CONFIG_DIRS=${XDG_CONFIG_DIRS:-/etc/xdg}:/usr/share/wuzetka-settings/xdg