# config/files/sneexy/mazurek/usr/etc/systemd/system/user@.service.d/local.conf
[Service]
Environment="PATH=/usr/lib/ccache/bin:/usr/local/sbin:/usr/local/bin:/usr/bin"
Environment="XDG_STATE_HOME=%h/.local/var/state"
